# supervisor/micha.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

; ── MICHA API (gunicorn) ─────────────────────────────────────────────────────
[program:micha_api]
command=/app/venv/bin/gunicorn config.wsgi:application -c config/gunicorn.conf.py
directory=/app
user=micha
autostart=true
autorestart=true               ; FIX: Auto-restart on crash
startretries=3
stdout_logfile=/var/log/micha/api.log
stderr_logfile=/var/log/micha/api_error.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=DJANGO_SETTINGS_MODULE="config.settings"
stopasgroup=true
killasgroup=true

; ── Celery Worker ─────────────────────────────────────────────────────────────
[program:micha_celery]
command=/app/venv/bin/celery -A config worker -l info -c 4 -Q high,default,low
directory=/app
user=micha
autostart=true
autorestart=true               ; FIX: Auto-restart if workers die
startretries=5
stdout_logfile=/var/log/micha/celery.log
stderr_logfile=/var/log/micha/celery_error.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stopasgroup=true
killasgroup=true
stopwaitsecs=60                ; Give workers 60s to finish current task on shutdown

; ── Celery Beat ───────────────────────────────────────────────────────────────
; FIX: Beat runs under supervisor — auto-restarts if it crashes
; Without this, all scheduled tasks (earnings release, price alerts etc.) stop silently
[program:micha_celery_beat]
command=/app/venv/bin/celery -A config beat -l info
directory=/app
user=micha
autostart=true
autorestart=true               ; FIX: Critical — Beat must never be down
startretries=5
stdout_logfile=/var/log/micha/celery_beat.log
stderr_logfile=/var/log/micha/celery_beat_error.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopasgroup=true
killasgroup=true

; ── Daphne (WebSocket/ASGI) ───────────────────────────────────────────────────
[program:micha_daphne]
command=/app/venv/bin/daphne -b 127.0.0.1 -p 8001 config.asgi:application
directory=/app
user=micha
autostart=true
autorestart=true
startretries=3
stdout_logfile=/var/log/micha/daphne.log
stderr_logfile=/var/log/micha/daphne_error.log
stopasgroup=true
killasgroup=true

[group:micha]
programs=micha_api,micha_celery,micha_celery_beat,micha_daphne
